# scripts/com.rstriq.portal-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  RSTR IQ — Verified Athletics portal auto-pull
  Fires every 4 hours at :00 (6am, 10am, 2pm, 6pm, 10pm, 2am).
  Install: see docs/PORTAL_AUTO_PULL.md
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rstriq.portal-pull</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>$HOME/dev-main/diamond-predictor-66/scripts/auto_pull_portal.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/rstr-iq-portal-pull.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/rstr-iq-portal-pull.err</string>
</dict>
</plist>
